# scripts/crontab
0 0 * * * /root/scripts/auto-fetch-metadata.sh
